# templates/default/ubuntu/upstart.freegeoip.conf.erb
description "freegeoip"
author "DigiTar"

start on (filesystem and net-device-up)
stop on runlevel [!2345]

respawn
respawn limit 5 30

chdir /usr/share/statsd
<% unless node["platform_version"].to_f < 11.10 -%>
setuid <%= node[:freegeoip][:user] %>
<% end %>

script
<% if node["platform_version"].to_f < 11.10 -%>
  exec start-stop-daemon --start --chuid <%= node[:freegeoip][:user] %> --exec <%= node[:freegeoip][:install_path]-%>/freegeoip -- -addr="<%=node[:freegeoip][:listen][:address]-%>:<%=node[:freegeoip][:listen][:port]-%>" <%= node[:freegeoip][:use_xforwardedfor] == true ? "-use-x-forwarded-for=true" : "" -%> >> <%= node[:freegeoip][:log_path] -%>/freegeoip.log 2>&1
<% else -%>
  exec <%= node[:freegeoip][:install_path]-%>/freegeoip -addr="<%=node[:freegeoip][:listen][:address]-%>:<%=node[:freegeoip][:listen][:port]-%>" <%= node[:freegeoip][:use_xforwardedfor] == true ? "-use-x-forwarded-for=true" : "" -%> >> <%= node[:freegeoip][:log_path] -%>/freegeoip.log 2>&1
<% end -%>
end script

emits statsd-running
